# misc/smd-pull.timer
[Unit]
Description=run smd-pull regularly

[Timer]
OnCalendar=*:0/2

[Install]
WantedBy=timers.target
